# deploy/systemd/agentic-sync-noche.timer
[Unit]
Description=Agentic Workflow — timer sync nocturno (10:45 UTC diario)

[Timer]
OnCalendar=*-*-* 10:45:00
# Si el servidor estaba apagado, ejecuta al arrancar si se perdió la ventana
Persistent=true

[Install]
WantedBy=timers.target
