# deploy/systemd/agentic-sync-mensual.service
[Unit]
Description=Agentic Workflow — sync mensual (cruceros + geo estado)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=alvaro.salis
WorkingDirectory=/home/alvaro.salis/agentic-workflow
Environment=PREFECT_API_URL=http://127.0.0.1:4200/api
ExecStart=/home/alvaro.salis/agentic-workflow/venv/bin/python scripts/sync_mensual.py
StandardOutput=journal
StandardError=journal
