# deploy/systemd/agentic-sync-noche.timer
[Unit]
Description=Agentic Workflow — timer sync nocturno (02:00 diario)

[Timer]
OnCalendar=*-*-* 02:00:00
# Si el servidor estaba apagado, ejecuta al arrancar si se perdió la ventana
Persistent=true

[Install]
WantedBy=timers.target
